# conf/supervisord.conf
[unix_http_server]
file = /opt/webmailbox/var/supervisor/supervisor.sock

[supervisord]
logfile = /opt/webmailbox/var/supervisor/supervisord.log
logfile_maxbytes = 50MB       ; maximum size of logfile before rotation
logfile_backups = 10          ; number of backed up logfiles
loglevel = info               ; info, debug, warn, trace
pidfile = /opt/webmailbox/var/supervisor/supervisord.pid
nodaemon = false              ; run supervisord as a daemon
minfds = 1024                 ; number of startup file descriptors
minprocs = 200                ; number of process descriptors
user = nobody                 ; default user
childlogdir = /opt/webmailbox/var/supervisor/log/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///opt/webmailbox/var/supervisor/supervisor.sock

[program:webmc-mongodb]
command = /opt/mongodb/bin/mongod -f /opt/webmailbox/conf/mongodb.conf

[program:webmc]
command = /opt/py26/bin/python /opt/webmailbox/application.py --port=8000
